# files/services.d/cron/run
#!/bin/execlineb -P
/usr/sbin/crond -f -l 8
